# root/etc/services.d/parsoid/run
#!/usr/bin/with-contenv bash
cd $PARSOID_HOME
su -c 'node bin/server.js' $PARSOID_USER